# packaging/systemd/olcrtc-manager.service
[Unit]
Description=LibreRTC Node
Documentation=https://github.com/svllvsxprod/librertc-node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=OLCRTC_PATH=/usr/local/bin/olcrtc
ExecStart=/usr/local/bin/librertc-node -config /etc/librertc/config.json
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=10s

# LibreRTC Node creates network namespaces, veth interfaces, tc qdiscs and iptables rules.
# Run as root unless you replace this with an equivalent capability setup.

[Install]
WantedBy=multi-user.target
